# etc/cron.d/ssh-geoip-alert
* * * * * root /usr/local/bin/ssh-geoip-alert
